# syslog-rx.socket
[Unit]
PartOf=syslog-rx.service

[Socket]
ListenStream=514
ListenDatagram=514

[Install]
WantedBy=sockets.target
